# notify_demo.service
[Unit]
Description=Notify Demo

[Service]
Type=notify-reload
NotifyAccess=all
ExecStart=/tmp/notify_demo

